# board/cuprous/gw/rootfs_overlay/etc/systemd/system/wpa_supplicant@.service.d/override.conf
[Service]
ExecStart=
ExecStart=/sbin/wpa_supplicant -g/var/run/wpa_supplicant -c/etc/wpa_supplicant/wpa_supplicant-%I.conf -i%I

